# .config/systemd/user/hyde-update-check.timer
[Unit]
Description=Weekly HyDE update check (Fridays)

[Timer]
OnCalendar=Fri *-*-* 12:00:00
Persistent=true

[Install]
WantedBy=timers.target
