# .config/systemd/user/noux-skills.service
[Unit]
Description=Noux Skills MCP Server
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/warley/life/2-areas/noux-skills
ExecStart=/usr/bin/node /home/warley/life/2-areas/noux-skills/dist/index.js
Restart=on-failure
RestartSec=5
Environment=NODE_ENV=production
Environment=NODE_TLS_REJECT_UNAUTHORIZED=0
Environment=PATH=/home/warley/.local/bin:/home/warley/.cargo/bin:/usr/local/bin:/usr/bin:/usr/sbin

[Install]
WantedBy=default.target
